# config/aws/moviemasher.upstart
start on stopped rc RUNLEVEL=[S3]
task
script
 cd /home/ec2-user/moviemasher.rb
 rake moviemasher:init >> /tmp/moviemasher/log/moviemasher.rb.init.log 2>&1 &
end script