# config/aws/moviemasher.upstart
start on stopped rc RUNLEVEL=[S3]
task
script
 cd /var/www/moviemasher.rb
 ./config/aws/create_temp_dirs.sh
 rake moviemasher:init['aws'] >> /tmp/moviemasher/log/moviemasher.rb.init.log 2>&1 &
end script
