# config/aws/moviemasher.cron
* * * * * cd /var/www/moviemasher.rb; export PATH=/usr/local/bin:/bin:/usr/bin:/usr/local/sbin:/usr/sbin:/sbin:/opt/aws/bin; rake moviemasher:process_queues >> /var/www/log/moviemasher.rb.cron.log 2>&1 &
