# config/moviemasher.cron

* * * * * cd /home/ec2-user/moviemasher.rb; /home/ec2-user/bin/rake moviemasher:process_queues >> /var/log/moviemasher/transcoder_cron.log 2>&1 &
