# config/azure/moviemasher.upstart
start on stopped rc RUNLEVEL=[S3]
task
script
 mkdir -p /tmp/moviemasher/queue
 mkdir -p /tmp/moviemasher/temporary
 mkdir -p /tmp/moviemasher/log
 mkdir -p /tmp/moviemasher/download
 mkdir -p /tmp/moviemasher/render
 chown -R azureuser:www-data /tmp/moviemasher
 chmod -R g+w /tmp/moviemasher/
 cd /home/azureuser/moviemasher.rb
 rake moviemasher:init['azure'] >> /tmp/moviemasher/log/moviemasher.rb.init.log 2>&1 &
end script
