# config/azure/moviemasher.cron
@reboot cd /var/www/moviemasher.rb; export PATH=/usr/local/bin:/bin:/usr/bin:/usr/local/sbin:/usr/sbin:/sbin; rake moviemasher:init['azure'] >> /var/log/moviemasher/moviemasher.rb.init.log 2>&1 &

* * * * * cd /var/www/moviemasher.rb; export TERM=dumb; export PATH=/usr/local/bin:/bin:/usr/bin:/usr/local/sbin:/usr/sbin:/sbin; rake moviemasher:process_queues >> /var/log/moviemasher/moviemasher.rb.cron.log 2>&1 &
